# cron/schedule.cron
# 0 * * * * $HOME/GitHub/geocoug/random-image/.venv/bin/python $HOME/GitHub/geocoug/random-image/random_image.py $HOME/GitHub/geocoug/random-image/tracker.json $HOME/iCloud/Pictures/Wallpaper
# 3 9 * * * $HOME/venvs/dev/bin/python $HOME/GitHub/geocoug/bitcoin/btc_tracker.py
1 9 * * * $HOME/venvs/dev/bin/python $HOME/GitHub/geocoug/nasa-apod-dl/nasa_apod_dl.py $HOME/GitHub/nasa.txt "$HOME/Library/Mobile Documents/com~apple~CloudDocs/Pictures/Wallpaper" > /dev/null 2>&1
# 9 9 * * * $HOME/bin/pg_backup.sh >> $HOME/iCloud/Backups/postgres/pg_backup.log 2>&1
13 9 * * 5 git --git-dir=$HOME/.oh-my-zsh/.git --work-tree=$HOME/.oh-my-zsh pull > /dev/null 2>&1 && git --git-dir=$HOME/.oh-my-zsh/custom/plugins/zsh-autosuggestions/.git --work-tree=$HOME/.oh-my-zsh/custom/plugins/zsh-autosuggestions pull > /dev/null 2>&1 && git --git-dir=$HOME/.oh-my-zsh/custom/plugins/zsh-syntax-highlighting/.git --work-tree=$HOME/.oh-my-zsh/custom/plugins/zsh-syntax-highlighting pull > /dev/null 2>&1
21 9 * * 5 cd $HOME/GitHub/geocoug/dotfiles/homebrew && brew bundle dump -f > /dev/null 2>&1
22 9 * * 5 crontab -l > $HOME/GitHub/geocoug/dotfiles/cron/schedule.cron
23 9 * * 5 tldr --update > /dev/null 2>&1
26 9 * * 5 cd $HOME/GitHub/geocoug/dotfiles && git add --all && git commit --message "backup $(date +'\%Y-\%m-\%d \%H:\%M:\%S')" > /dev/null 2>&1 && git push origin main > /dev/null 2>&1
